# Old_Versions/November/CMP9767M-master/catkin_ws/devel/etc/catkin/profile.d/90-uol_cmp9767m_base.sh
# extend the Gazebo model path with our models
export GAZEBO_MODEL_PATH=/home/computing/November/CMP9767M-master/catkin_ws/src/models:/home/computing/November/CMP9767M-master/catkin_ws/install/share/uol_cmp9767m_base/models:$GAZEBO_MODEL_PATH
